# home/dot_config/systemd/user/break-reminder.timer
[Unit]
Description=Break Reminder Timer

[Timer]
OnBootSec=1hour
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target
