# home/dot_config/systemd/user/reminder.timer
[Unit]
Description=Break Reminder

[Timer]
OnActiveSec=1hour
OnUnitInactiveSec=1hour

[Install]
WantedBy=timers.target
